# WiggleR/Code/Api/wiggler/service/wiggler.service
[Unit]
Description=WiggleR service

[Service]
WorkingDirectory=%h
ExecStart=/bin/bash /usr/bin/wiggler_boot.sh

[Install]
WantedBy=default.target